# infrastructure/spark/hd-worker-spark.service
[Unit]
Description=HD Worker (Supabase Job Processor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sparkuser
WorkingDirectory=/home/sparkuser/srv/hd-worker
EnvironmentFile=/home/sparkuser/srv/hd-worker/.env
ExecStart=/home/sparkuser/srv/hd-worker/.venv/bin/python3 /home/sparkuser/srv/hd-worker/hd_worker_mvp.py --loop --sleep 5
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
